# rnnoise-watch-default.service
[Unit]
Description=Watch default raw mic and relink it into RNNoise
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/.local/bin/rnnoise-watch-default.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
